# debian/postinst
#!/bin/bash
echo "reboot your system to start using joypad"
echo "or run systemctl start joypad.service"
systemctl enable joypad.service