# debian/postinst
#!/bin/bash

systemctl enable --now joypad.service