# debian/prerm
#!/bin/bash

systemctl disable --now joypad.service