# debian/prerm
#!/bin/bash

systemctl disable joypad.service